# apiAirlux/docker/conf/nginx.conf
[program:nginx]
command=/etc/nginx
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-err.log