# apiAirlux/docker/conf/ssh.conf
[program:ssh]
command=/usr/bin
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ssh.log
stderr_logfile=/var/log/supervisor/ssh-err.log